# target/linux/alpine/base-files/etc/uci-defaults/xx_customizations

wifi_key=$(uci -q get wireless.@wifi-iface[1].key)
if [ -z "$wifi_key" ]; then
	uci set wireless.@wifi-iface[1].ssid='home_network'
	uci set wireless.@wifi-iface[1].encryption=psk2
	uci set wireless.@wifi-iface[1].key='secret'

	uci set wireless.@wifi-device[1].country=US
	uci set wireless.@wifi-device[1].channel=auto
	uci set wireless.@wifi-device[1].disabled='0'

	uci commit wireless
fi

#uci set system.@system[0].hostname='OpenWrt-dev'
#uci commit system

#uci set network.lan.ipaddr=192.168.2.1
#uci commit network
